# ddns/crontab
*/10 * * * * python /usr/src/ddns/check.py

